# templates/etc/systemd/system/data-container.service.erb
[Unit]
Description=Docker Data Container for <%= @title %>
Requires=docker.service<% if @depends %><% @depends.each do |dep| %><%= " #{dep}.service" %><% end %><% end %>
After=docker.service<% if @depends %><% @depends.each do |dep| %><%= " #{dep}.service" %><% end %><% end %>

[Service]
Type=oneshot
Restart=no
RemainAfterExit=yes
Environment="<%= @title %>_DOCKER_IMAGE=<%= @image %>"
Environment="<%= @title %>_DOCKER_RUN_OPTS=<%= @docker_run_options %>"
ExecStart=-/usr/bin/docker run \
    $<%= @title %>_DOCKER_RUN_OPTS \
    $<%= @title %>_DOCKER_IMAGE
ExecStop=/usr/bin/docker stop <%= @title %>

[Install]
WantedBy=multi-user.target
